# config/systemd/user/ci_cmake.timer
[Unit]
Description=Hourly CMake CI Build

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
