# config/systemd/user/mspdbsrv.service
[Unit]
Description=mspdbsrv

[Service]
ExecStart=/usr/bin/sh %h/.config/msvc/mspdbsrv.sh

[Install]
WantedBy=default.target
