# config/systemd/user/newsbeuter.timer
[Unit]
Description=Hourly clean and update newsbeuter feeds

[Timer]
OnBootSec=1min
OnCalendar=hourly

[Install]
WantedBy=timers.target
