# systemd/user/checknews.timer
[Unit]
Description=Hourly clean and update news feeds

[Timer]
OnBootSec=1min
OnCalendar=hourly

[Install]
WantedBy=timers.target
